# systemd/fcgiwrap@.service.in
[Unit]
Description=Simple CGI Server for user %I
After=nss-user-lookup.target

[Service]
ExecStart=@sbindir@/fcgiwrap -s /run/fcgiwrap-%I.sock
User=%I
