# systemd/fcgiwrap@.socket.in
[Unit]
Description=fcgiwrap socket for user %I

[Socket]
ListenStream=/run/fcgiwrap-%I.sock

[Install]
WantedBy=sockets.target
